# systemd/yuri.service
[Unit]
Description=Yuri — @MyProjectWorldbot Telegram daemon (receipts photo intake + reports)
Documentation=https://github.com/galehomedesigns/openclaw-xrt9
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tonygale/openclaw
ExecStart=/home/tonygale/openclaw/.venv/bin/python /home/tonygale/openclaw/skills/receipts/scripts/yuri_daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yuri

[Install]
WantedBy=default.target
